# wps/usr/conf/s6/php/finish
#!/bin/sh

s6-svc -Dd /wps/run/php
rm -rf /wps/run/php