# wps/usr/conf/s6/php/run
#!/bin/sh

exec /usr/bin/php-fpm -c /home/wordpress/conf/php/php.ini -y /home/wordpress/conf/php/php-fpm.conf
