# wp-submarine/wps/usr/conf/s6/cron/run
#!/bin/sh

exec /usr/sbin/crond -f -L /var/log/cron.log
