# wps/usr/conf/s6/nginx/finish
#!/bin/sh

s6-svc -Dd $conf/init.d/nginx