# wps/usr/conf/s6/nginx/finish
#!/bin/sh

/usr/bin/s6-svscanctl -t $init