# wps/usr/conf/s6/nginx/finish
#!/bin/sh

s6-svc -Dd /wps/run/nginx
rm -rf /wps/run/nginx