# wps/usr/conf/s6/cron/finish
#!/bin/sh

s6-svc -Dd /wps/run/cron
rm -rf /wps/run/cron