# wps/usr/conf/s6/mysql/finish
#!/bin/sh

s6-svc -Dd /wps/run/mysql
rm -rf /wps/run/mysql