# wps/usr/conf/s6/mysql/finish
#!/bin/sh

/usr/bin/s6-svscanctl -t $conf/s6;