# wps/usr/conf/s6/mysql/run
#!/bin/sh

if [[  -d /var/lib/mysql  ]]; then
	exec /usr/bin/mysqld --user=mysql --datadir=/var/lib/mysql --log-error=/root/mysql.log
fi
