# wps/usr/conf/s6/.s6-svscan/finish
#!/bin/sh

case $@ in
	reboot)	echo "Reboot!"
	;;
esac

sudo rm -rf /wps/run
sudo kill -9 -1

if [[  $WPS_RESTART == 'true'  ]]; then
	echo "Restart!"
	exec sudo /usr/bin/wps start
fi

exit 0