# wps/usr/conf/s6/.s6-svscan/finish
#!/bin/sh

sudo kill -9 -1

if [[  -f /home/wordpress/.reboot  ]];
then su -l root -c "/usr/bin/wps start"
fi

exit 0