# wps/usr/conf/s6/.s6-svscan/finish
#!/bin/sh

sudo s6-svscanctl -st /service
sudo rm -rf /service
